# numa/rootfs/etc/s6-overlay/s6-rc.d/numa/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Numa add-on: numa
# Runs the Numa DNS resolver.
# RUST_LOG and config are prepared by init-numa before this service starts.
# ==============================================================================
bashio::log.info "Launching Numa..."
exec /usr/local/bin/numa /config/numa.toml
