# numa/rootfs/etc/s6-overlay/s6-rc.d/numa/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Numa
# s6-overlay service run script
# Starts the Numa DNS server with persistent config support.
#
# Numa CLI reference (confirmed from src/main.rs v0.13.1):
#   numa [config-path]   — config path is a POSITIONAL argument, not a flag.
#                          Defaults to numa.toml in the current directory.
#   GET /health          — returns {"status":"ok",...} when up (src/health.rs)
#
# Log level is controlled via the RUST_LOG environment variable (env_logger).
# ==============================================================================

declare numa_config
declare log_level

# ------------------------------------------------------------------------------
# Read add-on options
# ------------------------------------------------------------------------------
log_level=$(bashio::config 'log_level' 'info')
bashio::log.level "${log_level}"

# Map the HA log level to a RUST_LOG value so Numa's env_logger picks it up.
# Numa uses the standard env_logger crate — RUST_LOG is the only supported
# mechanism for controlling its internal log verbosity.
export RUST_LOG="${log_level}"

bashio::log.info "Starting Numa..."

# ------------------------------------------------------------------------------
# Config file provisioning
# ------------------------------------------------------------------------------
numa_config="/config/numa.toml"

if ! bashio::fs.file_exists "${numa_config}"; then
    bashio::log.info "No numa.toml found at ${numa_config}, writing default config..."
    cp /etc/numa/numa.toml.default "${numa_config}"
    bashio::log.info "Default config written. Edit ${numa_config} and restart the add-on to customise."
fi

bashio::log.info "Config : ${numa_config}"
bashio::log.info "Binary : $(numa version 2>&1 || true)"

# ------------------------------------------------------------------------------
# Launch Numa
#
# The config path is passed as a positional argument — NOT as --config.
# All DNS/API/port settings come from the TOML file itself; we do not
# patch them here: numa.toml is the single source of truth.
# ------------------------------------------------------------------------------
bashio::log.info "Launching Numa..."

exec /usr/local/bin/numa "${numa_config}"
