# numa/rootfs/etc/s6-overlay/s6-rc.d/numa/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: Numa
# s6-overlay service run script
# Starts the Numa DNS server with persistent config support.
#
# Numa CLI reference (confirmed from src/main.rs v0.13.1):
#   numa [config-path]   — config path is a POSITIONAL argument, not a flag.
#                          Defaults to numa.toml in the current directory.
#   GET /health          — returns {"status":"ok",...} when up (src/health.rs)
#
# Log level is controlled via the RUST_LOG environment variable (env_logger).
# ==============================================================================

declare numa_config
declare log_level
declare rust_log_level
declare config_source

# ------------------------------------------------------------------------------
# Read add-on options
# ------------------------------------------------------------------------------
log_level=$(bashio::config 'log_level' 'info')
bashio::log.level "${log_level}"

# Map the HA log level to a valid env_logger / RUST_LOG level.
# HA schema allows: trace | debug | info | notice | warning | error | fatal
# env_logger accepts: trace | debug | info | warn | error  (no notice/warning/fatal)
case "${log_level}" in
  notice)  rust_log_level="info"  ;;
  warning) rust_log_level="warn"  ;;
  fatal)   rust_log_level="error" ;;
  *)       rust_log_level="${log_level}" ;;
esac
export RUST_LOG="${rust_log_level}"

bashio::log.info "Starting Numa..."

# ------------------------------------------------------------------------------
# Config file provisioning
# ------------------------------------------------------------------------------
numa_config="/config/numa.toml"
config_source=$(bashio::config 'config_source' 'default')

if [ "${config_source}" = "file" ]; then
    # Strict mode: the user is responsible for providing the config file.
    if ! bashio::fs.file_exists "${numa_config}"; then
        bashio::log.fatal "config_source is set to 'file' but ${numa_config} does not exist."
        bashio::log.fatal "Create the file or change config_source to 'default' to auto-generate it."
        exit 1
    fi
    bashio::log.info "config_source=file: using existing ${numa_config}"
else
    # Default mode: copy the bundled template on first boot.
    if ! bashio::fs.file_exists "${numa_config}"; then
        bashio::log.info "No numa.toml found at ${numa_config}, writing default config..."
        cp /etc/numa/numa.toml.default "${numa_config}"
        bashio::log.info "Default config written. Edit ${numa_config} and restart the add-on to customise."
    fi
fi

bashio::log.info "Config : ${numa_config}"
bashio::log.info "Binary : $(numa version 2>&1 || true)"
# ------------------------------------------------------------------------------
# Network readiness check — two phases + circular-DNS workaround
#
# Numa's blocklist downloader fires immediately on startup (tokio::spawn with
# no delay). reqwest is compiled WITHOUT the hickory-dns feature (see upstream
# Cargo.toml), so it resolves hostnames via the libc getaddrinfo() syscall,
# which reads /etc/resolv.conf and queries whatever nameserver is listed there.
#
# On a typical home-network setup the router's DHCP server is configured to
# hand out the Home Assistant IP as the DNS server for LAN clients.  The
# add-on uses host_network: true, so the container inherits /etc/resolv.conf
# from the host — which therefore also points to Numa itself.  This creates a
# circular dependency at startup:
#
#   1. Numa bind()s port 53 (serve.rs, before the recv_from loop starts).
#   2. tokio::spawn(load_blocklists) fires on a worker thread immediately.
#   3. reqwest → getaddrinfo("cdn.jsdelivr.net") → DNS query → <HA_IP>:53.
#   4. Numa has the socket but the recv_from loop hasn't started yet.
#   5. The query sits unanswered until getaddrinfo times out → download fails.
#
# Circular-DNS fix — /etc/resolv.conf injection (bootstrap only, instant)
#   Prepend "nameserver 9.9.9.9" to the container's /etc/resolv.conf as a
#   TEMPORARY bootstrap measure.  Done immediately — no network access required.
#   The original file contents are saved to /tmp/resolv.conf.numa.orig and are
#   restored just before exec-ing Numa, so that steady-state DNS resolution goes
#   through Numa (which is running by then and can answer queries itself).
#
#   The injection is done in-place (cat > /etc/resolv.conf) rather than with mv,
#   so the original inode is preserved.  Docker bind-mounts resolv.conf into the
#   container; replacing the inode with mv would silently break that bind-mount
#   for the lifetime of the container.
#
#   Docker gives every container its own mount namespace, so this edit is local
#   to the container and does NOT touch the host's resolver or affect DHCP
#   clients on the LAN (they still use Numa as their DNS server).
#
# Network readiness check (brief, non-blocking)
#   Probes Quad9 by its hardcoded IP so this check cannot itself be blocked by a
#   broken system resolver.  Capped at 3 × 2 s = ~10 s worst case so the total
#   startup delay is bounded well within the HA watchdog window.  If the network
#   is not reachable yet we proceed immediately — Numa starts and the blocklist
#   will be retried after refresh_hours.
#
# NOTE: the proper long-term fix is upstream in Numa: download_blocklists()
# should use Numa's own configured DoH upstreams (or a hardcoded fallback IP)
# instead of relying on getaddrinfo().  A GitHub issue has been filed:
# https://github.com/razvandimescu/numa/issues  (see CHANGELOG for link).
# ------------------------------------------------------------------------------

# --- Circular-DNS fix: inject Quad9 into /etc/resolv.conf (bootstrap only) ---
# Done immediately — no network access required (just cp + cat).
# Skip silently if 9.9.9.9 is already present or the file is read-only.
# The original file is saved and will be restored just before exec.
declare _restore_resolv_conf=false
if ! grep -q '^nameserver 9\.9\.9\.9' /etc/resolv.conf 2>/dev/null; then
    bashio::log.info "Injecting Quad9 bootstrap resolver into /etc/resolv.conf (circular DNS fix)..."
    if cp /etc/resolv.conf /tmp/resolv.conf.numa.orig 2>/dev/null; then
        # Build new content in a temp file, then overwrite in-place so the
        # original inode is preserved (required for Docker bind-mount compatibility).
        { printf 'nameserver 9.9.9.9\n'; cat /tmp/resolv.conf.numa.orig; } > /tmp/resolv.conf.numa \
            && cat /tmp/resolv.conf.numa > /etc/resolv.conf \
            && _restore_resolv_conf=true \
            || bashio::log.warning "/etc/resolv.conf is read-only — circular DNS dependency may persist; blocklist download may fail."
    else
        bashio::log.warning "Could not back up /etc/resolv.conf — skipping Quad9 injection."
    fi
else
    bashio::log.debug "Quad9 already present in /etc/resolv.conf, no injection needed."
fi

# --- Network readiness check (brief, non-blocking) ---
# Probes Quad9 by hardcoded IP — cannot be blocked by a broken system resolver.
# Capped at 3 attempts × 2 s sleep = ~10 s worst case so Numa is always exec'd
# well within the HA watchdog / Docker HEALTHCHECK start-period window.
# On failure we proceed immediately; the blocklist will retry after refresh_hours.
declare _net_ready=false
for _attempt in 1 2 3; do
    if curl --silent \
            --max-time 2 --connect-timeout 2 \
            --output /dev/null \
            --resolve "dns.quad9.net:443:9.9.9.9" \
            "https://dns.quad9.net/dns-query?name=dns.quad9.net&type=A" 2>/dev/null; then
        _net_ready=true
        bashio::log.debug "Outbound network ready (attempt ${_attempt}/3)"
        break
    fi
    bashio::log.debug "Network not ready yet (attempt ${_attempt}/3)..."
    [ "${_attempt}" -lt 3 ] && sleep 2
done

if [ "${_net_ready}" != "true" ]; then
    bashio::log.warning "Outbound network not confirmed — proceeding anyway; blocklist will retry after refresh_hours."
fi

# ------------------------------------------------------------------------------
# Launch Numa
#
# The config path is passed as a positional argument — NOT as --config.
# All DNS/API/port settings come from the TOML file itself; we do not
# patch them here: numa.toml is the single source of truth.
# ------------------------------------------------------------------------------
# --- Restore original /etc/resolv.conf before handing off to Numa ---
# The bootstrap Quad9 entry is only needed during startup so that
# getaddrinfo() can reach the blocklist CDN before Numa's recv_from loop
# is ready.  Now that we are about to exec Numa, it will own port 53 and
# can answer DNS queries itself — so we restore the original file to
# preserve steady-state resolution order.
if [ "${_restore_resolv_conf}" = "true" ] && [ -f /tmp/resolv.conf.numa.orig ]; then
    cat /tmp/resolv.conf.numa.orig > /etc/resolv.conf \
        && bashio::log.debug "Restored original /etc/resolv.conf before exec." \
        || bashio::log.warning "Failed to restore /etc/resolv.conf — Quad9 will remain as bootstrap resolver."
fi

bashio::log.info "Launching Numa..."

exec /usr/local/bin/numa "${numa_config}"
